# infra/honeypot-cybermap-frontend.service
[Unit]
Description=Honeypot Cybermap Frontend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/honeypot-cybermap-frontend/dist
ExecStart=/usr/bin/python3 /opt/honeypot-cybermap-frontend/honeypot_static_server.py --host 0.0.0.0 --port 80 --directory /opt/honeypot-cybermap-frontend/dist --runtime-identity-file /opt/honeypot-cybermap-frontend/runtime_health_ed25519 --runtime-known-hosts-file /opt/honeypot-cybermap-frontend/known_hosts
Restart=always
RestartSec=3
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
