# infra/honeypot-flink-hdfs-sync.service
[Unit]
Description=Sync Flink honeypot stream output to HDFS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hadoop
Group=hadoop
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64
ExecStart=/bin/bash -lc '/opt/honeypot-bigdata/honeypot_flink_hdfs_sync.sh'
